# packaging/dure-server.service
[Unit]
Description=Dure central control plane
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/dure/server.env
ExecStart=/usr/bin/dure-server --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5
User=dure
Group=dure
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
